# ubuntu/etc/systemd/system/snapd.service.d/99-timeoutstopfix.conf
# https://www.omgubuntu.co.uk/2023/03/fix-ubuntu-slow-shutdowns-snap-daemon
# DOES NOT WORK
[Service]
TimeoutStopSec=5
